# build/resources/tsvauswertung.service
[Unit]
Description=TSV access Flask server
Wants=mariadb.target
After=mariadb.target

[Service]
Type=simple
WorkingDirectory=/opt/tsvserver
ExecStart=/usr/bin/python3 /opt/tsvserver/TsvAuswertung.py

[Install]
WantedBy=multi-user.target
