# build/resources/tsvaccess.service
[Unit]
Description=TSV access rfid deamon
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/TSVAccess
ExecStart=/usr/bin/python3 /opt/TSVAccess/TsvAccessModule.py

[Install]
WantedBy=multi-user.target

