# build/resources/tsvaccess.service
[Unit]
Description=TSV access rfid deamon
Wants=network-online.target
After=network-online.target
StartLimitIntervalSec=600
StartLimitBurst=3

[Service]
Restart=on-failure
RestartSec=5s
Type=simple
WorkingDirectory=/opt/TSVAccess
ExecStart=/usr/bin/python3 /opt/TSVAccess/TsvAccessModule.py
#use this for LEDs only (not relais)
#ExecStart=/usr/bin/python3 /opt/TSVAccess/TsvAccessModule.py -i

[Install]
WantedBy=multi-user.target

