# build/resources/tsvbackup.service
[Unit]
Description=Backup service

[Service]
ExecStart=/root/bin/backup.sh