# build/resources/tsvExternalDLDaily.service
[Unit]
Description=TSV external info service quarter hour current access

[Service]
ExecStart=/usr/bin/python3 /opt/tsvserver/ExternalMailService.py --dl_daily
 
