# build/resources/tsvExternalDLDaily.timer
[Unit]
Description=TSV external timer for 15 min current access

[Timer]
OnCalendar=*:0/15
Persistent=true
Unit=tsvExternalDLDaily.service

[Install]
WantedBy=timers.target
 
