# build/resources/tsvSewobe.service
[Unit]
Description=TSV Sewobe member update

[Service]
Type=oneshot
WorkingDirectory=/opt/tsvserver
ExecStart=/usr/bin/python3 /opt/tsvserver/SewobeConnector.py
 
