# build/resources/tsvSewobe.service
[Unit]
Description=TSV Sewobe member update

[Service]
ExecStart=/usr/bin/python3 /opt/tsvserver/ExternalMailService.py --sewobe
 
